# deploy/systemd/kk-backend.service
# Kuteera Kitchen — FastAPI backend service
# Place at: /etc/systemd/system/kk-backend.service
# Enable:   sudo systemctl enable kk-backend
# Start:    sudo systemctl start kk-backend

[Unit]
Description=Kuteera Kitchen FastAPI backend
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/kk_v1

# Load secrets from a dedicated env file — never commit real credentials
EnvironmentFile=/var/www/kk_v1/backend/.env

ExecStart=/var/www/kk_v1/.venv/bin/uvicorn backend.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2

Restart=on-failure
RestartSec=5s

# Harden: no write access outside WorkingDirectory
ReadWritePaths=/var/www/kk_v1
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
