# Docs/vps-setup/kk-frontend-prod.service
[Unit]
Description=Kuteera Kitchen Next.js Frontend (PROD)
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/var/www/kk_v1_prod
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kk-frontend-prod
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=BACKEND_URL=http://127.0.0.1:8000

[Install]
WantedBy=multi-user.target
